# .config/autostart/kanshi.desktop
[Desktop Entry]
Type=Application
Name=kanshi
Exec=kanshi
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
